# chezmoi/Library/LaunchAgents/com.macmikase.caps-to-f18.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.macmikase.caps-to-f18</string>

    <key>ProgramArguments</key>
    <array>
      <string>{{ .chezmoi.homeDir }}/.local/bin/macmikase-remap-caps-to-f18.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>{{ .chezmoi.homeDir }}/Library/Logs/com.macmikase.caps-to-f18.log</string>

    <key>StandardErrorPath</key>
    <string>{{ .chezmoi.homeDir }}/Library/Logs/com.macmikase.caps-to-f18.log</string>
  </dict>
</plist>
